# server/supervisor/development.ini
[program:harvest_project-development]
directory = /home/devuser/webapps/harvest_project-env/harvest_project
command = /home/devuser/webapps/harvest_project-env/bin/uwsgi --ini /home/devuser/webapps/harvest_project-env/harvest_project/server/uwsgi/development.ini
autostart = true
autorestart = true
