# server/supervisor/local.ini
[program:harvest_project-local]
directory = /home/devuser/.virtualenvs/harvest_project-env/harvest_project
command = /home/devuser/.virtualenvs/harvest_project-env/bin/uwsgi --ini /home/devuser/.virtualenvs/harvest_project-env/harvest_project/server/uwsgi/local.ini
autostart = true
autorestart = true
